# systemd/alexa.service
[Unit]
Description=Amazon Alexa Voice Service
Wants=network-online.target
After=network-online.target

[Service]
Type=forking
ExecStart=sudo /home/__USER__/Assistants-Pi/Alexa/start.sh

WorkingDirectory=/home/__USER__/
StandardOutput=inherit
StandardError=inherit
Restart=always
User=__USER__

[Install]
WantedBy=multi-user.target
